# systemd/vpnje-bot.service
[Unit]
Description=Vpnje bot service.

[Service]
Type=simple
ExecStart=/usr/bin/nodejs /root/vpnje-bot/app.js

[Install]
WantedBy=multi-user.target